# native/KaisolaMac/build/Build/Products/Debug/KaisolaMacPreview.app/Contents/Library/LaunchAgents/com.kaisola.mac.preview.broker-bootstrap.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.kaisola.mac.preview.broker-bootstrap</string>
    <key>BundleProgram</key>
    <string>Contents/Resources/BrokerHelper/bin/kaisola-broker-bootstrap</string>
    <key>MachServices</key>
    <dict>
        <key>com.kaisola.mac.preview.broker-bootstrap</key>
        <true/>
    </dict>
    <key>ProcessType</key>
    <string>Interactive</string>
    <key>RunAtLoad</key>
    <false/>
    <key>AssociatedBundleIdentifiers</key>
    <array>
        <string>com.kaisola.mac.preview</string>
    </array>
</dict>
</plist>
